# src/services/display/pioled_stats.service
[Unit]
Description=PiOLED stats display service

[Service]
Type=simple
User=niraj
ExecStart=/bin/sh -c "python3 /home/niraj/projects/laika/src/services/display/stats.py"
Restart=always

[Install]
WantedBy=multi-user.target